# recipes-core/base-files/base-files/media-sdcard.automount
[Unit]
Description=Automount for SDCard
ConditionPathExists=/dev/mmcblk1p1
RefuseManualStart=true
RefuseManualStop=true
After=local-fs.target

[Automount]
Where=/media/sdcard

[Install]
WantedBy=local-fs.target
