# recipes-support/rng-tools/files/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
ExecStart=/usr/sbin/rngd -f -r /dev/hwrng

[Install]
WantedBy=multi-user.target
