# recipes-support/aero-http-server/files/aero-http-server.socket
[Unit]
Description=aero-http-server socket
PartOf=aero-http-server.service

[Socket]
ListenStream=8000

[Install]
WantedBy=sockets.target
